# roles/druid.broker/templates/supervisord.conf
[program:broker]
directory = {{ druid_base_dir }}/druid
command = java
	-server
	-Xms22g
	-Xmx22g
    -XX:+UseConcMarkSweepGC
    -XX:MaxDirectMemorySize=6g
    -Duser.timezone=UTC
    -Dfile.encoding=UTF-8
    -Djava.util.logging.manager=org.apache.logging.log4j.jul.LogManager
    -Djava.io.tmpdir={{ druid_tmpdir }}
    -cp {{ druid_conf_dir }}/druid/_common:{{ druid_conf_dir }}/druid/broker:{{ druid_base_dir }}/druid/lib/*
    {% if druid_package_type == 'apache' %}{{ apache_druid_main_class }}{% else %}{{ io_druid_main_class }}{% endif %} server broker
user = {{ druid_run_user }}
environment =
    nodeType=broker,
    DRUID_HOME="{{ druid_base_dir }}/druid",
stdout_logfile = /var/log/druid/broker.log
stdout_logfile_backups = 0
redirect_stderr = true
autostart = true
autorestart = true
startretries=3